# deploy/systemd/aistock-api.service
[Unit]
Description=AiStock FastAPI
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/aistock/backend
EnvironmentFile=/opt/aistock/deploy/aistock.env
ExecStart=/opt/aistock/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
